# dot_config/systemd/user/udiskie.service
[Unit]
Description=Automated devices mounting tool

[Service]
Type=simple
ExecStart=/usr/bin/udiskie
Restart=on-failure

