# hosts/lumo/home/services/tempestmiku/openrc-service
#!/sbin/openrc-run
name="lumo-tempestmiku"
description="TempestMiku API and worker container"
supervisor=supervise-daemon
command="@runner@"
command_user="root"
output_log="/var/log/lumo/tempestmiku.log"
error_log="/var/log/lumo/tempestmiku.log"
respawn_delay=10
respawn_max=0

depend() {
  need lumo-podman lumo-postgresql lumo-cliproxyapi lumo-ntfy lumo-tempestmiku-embeddings
  after networking
}

start_pre() {
  checkpath -d -m 0750 -o tempestmiku:tempestmiku @dataDir@
  checkpath -f -m 0640 -o root:root /var/log/lumo/tempestmiku.log

  source_archive="@buildContext@/source.tar.gz"
  if [ ! -f "$source_archive" ] ||
    ! printf '%s  %s\n' '@sourceArchiveSha256@' "$source_archive" | sha256sum -c - >/dev/null; then
    eend 1 "TempestMiku source archive is missing or has the wrong digest"
    return 1
  fi

  push_key="$(tr -d '\r\n' < '@pushKeyPath@')"
  api_key="$(tr -d '\r\n' < '@sharedApiKeyPath@')"
  if [ -z "$push_key" ] || [ -z "$api_key" ]; then
    eend 1 "TempestMiku secrets are empty"
    return 1
  fi

  {
    printf '%s\n' 'HOME=/data'
    printf '%s\n' 'RUST_LOG=tm_server=info,tower_http=info'
    printf '%s\n' 'TM_SERVER_ADDR=127.0.0.1:@port@'
    printf '%s\n' 'TM_SERVER_ROLE=all'
    printf '%s\n' 'TM_DATABASE_URL=host=/run/postgresql port=@postgresPort@ user=tempestmiku dbname=tempestmiku'
    printf '%s\n' 'TM_AUTH_MODE=device'
    printf '%s\n' 'TM_PUBLIC_BASE_URL=@publicUrl@'
    printf '%s\n' 'TM_OMP_ACP_ENABLED=0'
    printf '%s\n' 'TM_OMP_ACP_ARTIFACT_ROOT=/data/artifacts'
    printf '%s\n' 'TM_MANAGED_SKILLS_PATH=/data/managed-skills'
    printf '%s\n' 'TM_MANAGED_MODE_ADDENDA_PATH=/data/managed-mode-addenda'
    printf '%s\n' 'TM_MANAGED_PERSONA_ADDENDA_PATH=/data/managed-persona-addenda'
    printf '%s\n' 'TM_MEMORY_EMBEDDING_PROVIDER=local'
    printf '%s\n' 'TM_MEMORY_EMBEDDING_ENDPOINT=http://127.0.0.1:11434/v1/embeddings'
    printf '%s\n' 'TM_MEMORY_EMBEDDING_MODEL=bge-m3:567m'
    printf '%s\n' 'TM_MEMORY_EMBEDDING_DIMENSIONS=1024'
    printf '%s\n' 'TM_MEMORY_EMBEDDING_NORMALIZATION=l2'
    printf '%s\n' 'TM_MEMORY_EMBEDDING_TIMEOUT_MS=10000'
    printf '%s\n' 'TM_MEMORY_EMBEDDING_MAX_BATCH_SIZE=16'
    printf '%s\n' 'TM_MEMORY_EMBEDDING_MAX_INPUT_BYTES=16384'
    printf '%s\n' 'TM_SELF_HOSTED_ASR_ENDPOINT=@selfHostedAsrEndpoint@'
    printf '%s\n' 'TM_SELF_HOSTED_ASR_LABEL=Home TEA-ASR (Taiwan Mandarin)'
    printf '%s\n' 'TM_SELF_HOSTED_ASR_MODEL_ID=JacobLinCool/TEA-ASR-1.1-mini'
    printf '%s\n' 'TM_PUSH_PROVIDER=unifiedpush'
    printf '%s\n' 'TM_UNIFIED_PUSH_ENDPOINT_ORIGIN=@pushOrigin@'
    printf 'TM_PUSH_ENCRYPTION_KEY=%s\n' "$push_key"
    printf '%s\n' 'OPENAI_BASE_URL=@openaiBaseUrl@'
    printf 'OPENAI_API_KEY=%s\n' "$api_key"
    printf '%s\n' 'OPENAI_MODEL=gpt-5.5'
    printf '%s\n' 'OPENAI_STREAM_USAGE=0'
  } > /run/lumo-tempestmiku.env
  chmod 0400 /run/lumo-tempestmiku.env

  if ! @podman@ image exists @image@; then
    @podman@ build \
      --build-arg SOURCE_REV=@sourceRev@ \
      --build-arg SOURCE_ARCHIVE_SHA256=@sourceArchiveSha256@ \
      --file @dockerfile@ \
      --network=host \
      --tag @image@ \
      @buildContext@ >&2
  fi
}
