# hosts/lumo/home/services/hermes-agent/openrc-service
#!/sbin/openrc-run
name="lumo-hermes-agent"
description="Lumo Hermes Agent (Nous Research) gateway container"
supervisor=supervise-daemon
command="@podman@"
command_args="run --replace --rm --network=host --env=HERMES_HOME=/opt/data --env=HERMES_DISABLE_LAZY_INSTALLS=0 --entrypoint=/bin/sh --name=lumo-hermes-agent -v @dataDir@:/opt/data -v @workDir@:@workDir@ -v @containerEntrypoint@:/usr/local/bin/hermes-honcho-entrypoint:ro @image@ /usr/local/bin/hermes-honcho-entrypoint"
command_user="root"
output_log="/var/log/lumo/hermes-agent.log"
error_log="/var/log/lumo/hermes-agent.log"
respawn_delay=10
respawn_max=0

depend() {
  need lumo-podman lumo-cliproxyapi lumo-honcho-api
  after networking
}

start_pre() {
  checkpath -f -m 0640 -o root:root /var/log/lumo/hermes-agent.log
  # Only pull when the image is missing. The ~3.85 GB image is re-copied on
  # every pull, so unconditional pulls would stall each service restart.
  if ! @podman@ image exists @image@; then
    @podman@ pull @image@ >&2
  fi
}
